# crontab
* * * * * node /root/src/index.js /root/data >> /var/log/cron.log 2>&1
